# deploy/flask-service.service
[Unit]
Description=My Pricer Flask API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=$DEPLOY_DIR/web/backend
EnvironmentFile=/etc/my-pricer/secrets/backend.env
ExecStart=/opt/crypto-pricer/venv/bin/python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
